# ecommerce.service
[Unit]
Description=E-Commerce API Docker Compose Application
Requires=docker.service
After=docker.service
Documentation=file:///opt/ecommerce/Multi-Vendor-E-Commerce-API/DEPLOYMENT_GUIDE.md

[Service]
Type=oneshot
WorkingDirectory=/opt/ecommerce/Multi-Vendor-E-Commerce-API
User=root
Group=root

# Start the application
ExecStart=/usr/local/bin/docker-compose -f docker-compose.prod.yml up -d

# Stop the application
ExecStop=/usr/local/bin/docker-compose -f docker-compose.prod.yml down

RemainAfterExit=yes

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=3

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security - Run as root (for Docker management)
# Consider using rootless Docker for better security in production
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
